# pipline/crontab
56 20 * * * /venv/bin/python3 /app/pipline.py > /app/cron_log.log 2>&1
